# jetsonclocks.service
[Unit]
Description=Run Jetson Clocks Script
Type=oneshot

[Service]
ExecStart=/usr/bin/jetson_clocks

[Install]
WantedBy=multi-user.target
